# fs.overlay/etc/hotplug.d/button/00-logger
#!/bin/sh
logger -t button "event: BUTTON=$BUTTON ACTION=$ACTION SEEN=$SEEN"
